# 05_services/scripts/launchd/com.sunshine.cover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sunshine.cover</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/al02399300/Desktop/Jun_AI/Dog_Contents/project_sunshine/services/scripts/cover_processor.py</string>
        <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/al02399300/Desktop/Jun_AI/Dog_Contents/project_sunshine</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/al02399300/Desktop/Jun_AI/Dog_Contents/project_sunshine/config/logs/cover_processor_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/al02399300/Desktop/Jun_AI/Dog_Contents/project_sunshine/config/logs/cover_processor_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/al02399300/Desktop/Jun_AI/Dog_Contents/project_sunshine</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>ExitTimeOut</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
